# com.github.essandess.adblock2privoxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.essandess.adblock2privoxy</string>
	<key>Program</key>
	<string>/bin/bash</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
 		<string>-c</string>
	<string>/bin/mkdir -p /usr/local/etc/adblock2privoxy &amp;&amp; ( /bin/test -f /usr/local/etc/adblock2privoxy/privoxy/ab2p.task &amp;&amp; /usr/local/bin/adblock2privoxy -t /usr/local/etc/adblock2privoxy/privoxy/ab2p.task || /usr/local/bin/adblock2privoxy -p /usr/local/etc/adblock2privoxy/privoxy -w /usr/local/etc/adblock2privoxy/css -d 127.0.0.1:8119 https://easylist.to/easylist/easyprivacy.txt https://easylist.to/easylist/easylist.txt https://easylist.to/easylist/fanboy-annoyance.txt https://easylist.to/easylist/fanboy-social.txt https://easylist-downloads.adblockplus.org/antiadblockfilters.txt https://easylist-downloads.adblockplus.org/malwaredomains_full.txt https://raw.githubusercontent.com/ryanbr/fanboy-adblock/master/fanboy-antifacebook.txt https://raw.githubusercontent.com/Dawsey21/Lists/master/adblock-list.txt ) &amp;&amp; /bin/launchctl unload -w /Library/LaunchDaemons/org.macports.Privoxy.plist &amp;&amp; /bin/launchctl load -w /Library/LaunchDaemons/org.macports.Privoxy.plist</string>
	</array>
        <key>RunAtLoad</key>
        <false/>
        <key>StartCalendarInterval</key>
        <array>
                <dict>
                        <key>Weekday</key>
                        <integer>7</integer>
                        <key>Hour</key>
                        <integer>1</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                </dict>
        </array>
	<key>StandardErrorPath</key>
	<string>/var/log/system.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/system.log</string>
</dict>
</plist>
